# system_files/etc/profile.d/gamersoulslop.sh
# shellcheck disable=SC2148
#SHUT THE FUCK UP SHELLCHECK
export PROTON_ENABLE_WOW64=1
export PROTON_USE_WOW64=1
export PROTON_ENABLE_NTSYNC=1
export PROTON_USE_NTSYNC=1
export PROTON_ENABLE_WAYLAND=1
export PROTON_USE_WAYLAND=1
export PROTON_PREFER_SDL=1


#Shoutouts to tulip for the banger gamersloppy 
https://github.com/zirconium-dev/zirconium-jackrabbit/blob/main/system_files/etc/profile.d/gamerslop.sh